# samples/OpenRC runscript/acpid
#!/sbin/openrc-run

description="Daemon for Advanced Configuration and Power Interface"

extra_started_commands="reload"
command="/usr/sbin/acpid"
command_args="$ACPID_ARGS"
start_stop_daemon_args="--quiet"

depend() {
	need localmount
	use logger
}

reload() {
	ebegin "Reloading acpid configuration"
	start-stop-daemon --exec $command --signal HUP
	eend $?
}
