# samples/desktop/nyancat.socket
[Unit]
Description=Nyan.cat socket

[Socket]
ListenStream=23
Accept=true

[Install]
WantedBy=sockets.target
